# .config/autostart/RGB.sh.desktop
[Desktop Entry]
Exec=/home/steinardth/Forrit/Script/RGB.sh
Icon=application-x-shellscript
Name=RGB.sh
Type=Application
X-KDE-AutostartScript=true
